# roles/qbittorrent/templates/setup-port-forward.timer.j2
# {{ ansible_managed }}
# vi ft=systemd

[Unit]
Description=Setup Port Forwarding for qBittorrent Timer

[Timer]
OnCalendar={{ port_forward_timer_schedule }}

[Install]
WantedBy=timers.target
